# etc/systemd/system/wazo-webhookd.service.d/create-dedicated-exchange.conf
[Service]
ExecStartPre=/usr/bin/wazo-webhookd-dedicated-exchange-create
